# docker/php/settings/supervisor-apache.conf
[supervisord]
nodaemon=true

[program:apache2]
command=apache2ctl -D "FOREGROUND" -k start
autostart=true
autorestart=true
redirect_stderr=true